# pkgs/workbench/cldemo-wbench-hortonworks-2lt22s-ansible/debian/home/cumulus/ansibledemos/roles/bigdata-common/templates/server-init-script.j2
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          {{item.serviceclass}}-{{item.service}}
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the {{item.serviceclass}} {{item.service}} service
# Description:       System mode startup script for
#                    the {{item.serviceclass}} {{item.service}} service.
### END INIT INFO

SERVICECLASS={{item.serviceclass}} # hadoop, yarn, mapred
USER={{item.user}} # hdfs, yarn, mapred
SERVICE={{item.service}} # namenode, datanode, resourcemanager,
                         # nodemanager, historyserver

# DO NOT EDIT ANYTHING BELOW

. /lib/lsb/init-functions

# XXX make these log to console if they fail
test -f /etc/hadoop/conf/directories.sh && \
    . /etc/hadoop/conf/directories.sh || exit 1
test -f /etc/hadoop/conf/usersAndGroups.sh && \
    . /etc/hadoop/conf/usersAndGroups.sh || exit 1
test -f /etc/hadoop/conf/$SERVICECLASS-env.sh && \
    USER=$USER . /etc/hadoop/conf/$SERVICECLASS-env.sh

SU_OPT="-l"

if [ $SERVICECLASS = "hadoop" ]; then
	PIDDIR=$HADOOP_PID_DIR
    DAEMON=/usr/lib/hadoop/sbin/hadoop-daemon.sh
elif [ $SERVICECLASS = "yarn" ]; then
	PIDDIR=$YARN_PID_DIR
    DAEMON=/usr/lib/hadoop-yarn/sbin/yarn-daemon.sh
elif [ $SERVICECLASS = "mapred" ]; then
	PIDDIR=$MAPRED_PID_DIR
    DAEMON=/usr/lib/hadoop-mapreduce/sbin/mr-jobhistory-daemon.sh
    SU_OPT=""
fi

PIDFILE=$PIDDIR/$SERVICECLASS-$USER-$SERVICE.pid
DAEMONUSER=$USER
DAEMONGROUP=$SERVICECLASS
PATH=/sbin:/bin:/usr/sbin:/usr/bin # XXX needed?

test -x $DAEMON || exit 0

if [ "$DAEMONUSER" = "$YARN_USER" ]; then
    DAEMON="export HADOOP_LIBEXEC_DIR=/usr/lib/hadoop/libexec; $DAEMON"
fi

do_start () {
	log_daemon_msg "Starting $SERVICECLASS-$SERVICE service"
	if [ ! -d $PIDDIR ]; then
		mkdir -p $PIDDIR
		chown $DAEMONUSER:$DAEMONGROUP $PIDDIR
	fi
	su $DAEMONUSER $SU_OPT -c "$DAEMON --config $HADOOP_CONF_DIR start $SERVICE &>/dev/null"
	status=$?
	log_end_msg ${status}
}

do_stop () {
	log_daemon_msg "Stopping $SERVICECLASS-$SERVICE service"
	su $DAEMONUSER $SU_OPT -c "$DAEMON --config $HADOOP_CONF_DIR stop $SERVICE &>/dev/null"
	log_end_msg $?
}

case "$1" in
	start|stop)
		do_${1}
		;;
	restart|reload|force-reload)
		if [ -s $PIDFILE ] && kill -0 $(cat $PIDFILE) >/dev/null 2>&1; then
			do_stop
			do_start
		fi
		;;
	force-stop)
		pulseaudio_stop
		killall pulseaudio || true
		sleep 2
		killall -9 pulseaudio || true
		;;
	status)
        if [ -s $PIDFILE ]; then
            pid=$(cat $PIDFILE)
            if ps --no-heading -o command $pid 2>/dev/null | awk '{print $NF}'|  grep -iq $SERVICE; then
                echo "$SERVICECLASS-$SERVICE running on pid $pid"
                exit 0
            else
                echo "$SERVICECLASS-$SERVICE stopped. Stale pid file."
                exit 1
            fi
        else
            echo "$SERVICECLASS-$SERVICE stopped."
            exit 3
        fi
		;;
	*)
		echo "Usage: /etc/init.d/$SERVICECLASS-$SERVICE {start|stop|force-stop|restart|reload|force-reload|status}"
		exit 1
		;;
esac

exit 0
